# mosaicme/web/mosaicme_supervisor.conf
[program:mosaicme_web]
command = gunicorn --name=mosaicme --bind=0.0.0.0:9000 --workers 3 web.wsgi:application
redirect_stderr = true
autostart = true
